# training-dev/workspace/boxes/chef/chef-repo/cookbooks/psdash/templates/default/psdash.erb
#!/bin/bash
#
# Startup script for psdash
#
# chkconfig: - 13 87
# description: psdash a dashboard for the psutils
# processname: psdash


# Source function library.
. /etc/rc.d/init.d/functions

LOGDIR="/var/log/psdash"
LOGFILE="${LOGDIR}/psdash.log"
PROG="/usr/bin/psdash"

start(){
  echo -n "Starting $PROG:"
  mkdir -p $LOGDIR
  nohup $PROG <%= @logs %> -b <%= @host %> -p <%= @port %> >$LOGFILE 2>&1 &
  RUNNING=`ps ax | grep -v grep | grep -c $PROG`
  if [ $RUNNING -gt 0 ]; then
    echo_success
  else
    echo_failure
  fi
  echo
}


stop(){
  PID=`ps x | grep -v grep | grep $PROG | awk '{print $1}'`;
  echo -n "Stopping $PROG:"
  if [ "x$PID" != "x" ]; then
    kill $PID;
  fi
  RUNNING=`ps ax | grep -v grep | grep -c $PROG`
  if [ $RUNNING -gt 0 ]; then
    echo_failure
  else
    echo_success
  fi
  echo
}

restart(){
  stop
  start
}

status(){
  if ps ax | grep -v grep | grep -c $PROG; then
    echo "psdash is running"
    exit 0
  else
    echo "psdash is not running"
    exit 3
  fi
}

case "$1" in
  start)
    start
  ;;
  stop)
    stop
  ;;
  restart)
    restart
  ;;
  status)
    status
  ;;
  *)
    echo $"Usage: $0 {start|stop|restart|status}"
    exit 2
esac
